# config/supervisord.conf
; supervisord runs as PID 1 in the container and manages two processes:
;   - pzserver: the Project Zomboid dedicated server (launched as the steam user)
;   - webui:    the Node.js admin UI (also steam, so config stays steam-owned)
; The unix socket is group-accessible by `steam` so the web UI can drive
; supervisorctl (start/stop/restart) without root or docker.sock.

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/run/supervisord.pid
childlogdir=/var/log/supervisor

[unix_http_server]
file=/run/supervisor.sock
chmod=0770
chown=root:steam

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[program:pzserver]
command=/server/scripts/launch_pz.sh
user=steam
directory=/home/steam/pz-dedicated
autostart=true
autorestart=unexpected
startsecs=15
stopsignal=TERM
stopwaitsecs=90
stdout_logfile=/var/log/supervisor/pzserver-stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=2
redirect_stderr=true
environment=HOME="/home/steam"

[program:webui]
command=node /server/webui/server.js
user=steam
directory=/server/webui
autostart=true
autorestart=true
startretries=5
stdout_logfile=/var/log/supervisor/webui-stdout.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
redirect_stderr=true
environment=HOME="/home/steam",NODE_ENV="production"
